# trunk/package/base-files/files/etc/inittab
::sysinit:/etc/rc.sysinit
@tty@::askfirst:/bin/sh
::ctrlaltdel:/sbin/reboot
